# 4.9.1/overlay/etc/update-motd.d/70-mattermost
#!/bin/bash

if [ ! -f /root/mattermost/done ]
then
cat <<EOF
We are installing mattermost, you can follow installation with the command below:

$> tail -f /var/log/mattermost.log

EOF
fi
