# overlay/etc/update-motd.d/70-mattermost
#!/bin/bash


cat << EOF
Mattermost
EOF
