# 4.9.1/overlay/etc/supervisor/conf.d/mattermost.conf
[program:mattermost-platform]
directory=/root/mattermost
command=/root/mattermost/bin/platform
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
